# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:llama-server]
command=llama-server -m %(ENV_HOME)s/app/models/ckpt-001.gguf --port 9000 --host 0.0.0.0 --ctx-size 2048 -t 2
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5

[program:api]
command=%(ENV_HOME)s/app/.venv/bin/uvicorn gazet.api:app --host 0.0.0.0 --port 8000
directory=%(ENV_HOME)s/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=10

[program:streamlit]
command=%(ENV_HOME)s/app/.venv/bin/streamlit run gazet_demo.py --server.port 7860 --server.address 0.0.0.0
directory=%(ENV_HOME)s/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=15
